# systemd/v0-farmar-monitor.timer
[Unit]
Description=v0 Farm Console - operational health monitor timer (criarov0)

[Timer]
# Frequent enough to catch an outage quickly without hammering the DB/web
# process; independent of any other project's timers on this host.
OnCalendar=*:0/5
RandomizedDelaySec=20
Persistent=true
Unit=v0-farmar-monitor.service

[Install]
WantedBy=timers.target
